# init.d/xremap
#!/sbin/openrc-run

command="/usr/bin/xremap"
command_args="--watch=device,config ~yamad/.config/xremap.yml"
command_background="true"
pidfile="/run/xremap.pid"

depend() {
	need localmount # modules(uinput)
}
